# daemon/com.upkeep.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.upkeep.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/usr/local/lib/upkeep/upkeep_daemon.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/upkeep-daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/upkeep-daemon.err</string>

    <key>WorkingDirectory</key>
    <string>/usr/local/lib/upkeep</string>

    <key>UserName</key>
    <string>root</string>

    <key>GroupName</key>
    <string>wheel</string>

    <key>ThrottleInterval</key>
    <integer>1</integer>
</dict>
</plist>
